# image/sling.service
[Unit]
Description=Sling
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
TimeoutStartSec=0
Restart=always
RestartSec=2
User=robot
ExecStart=/usr/local/bin/start-sling.sh

[Install]
WantedBy=multi-user.target
